# scripts/jebao-mqtt-update.timer
# /etc/systemd/system/jebao-mqtt-update.timer
#
# Timer to check for updates every 15 minutes

[Unit]
Description=Check for Jebao MQTT Bridge updates

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
